# .gemini/infra/systemd/sleep-cycle.timer
[Unit]
Description=Timer for Nightly GeminiClaw Sleep Cycle

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
