# .gemini/infra/systemd/tank-inspector.timer
[Unit]
Description=Run Vopak Tank Inspector hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
