# .gemini/harvester.timer
[Unit]
Description=Run Vopak Assistant Harvester daily at 8am

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
