# .agent/harvester.service
[Unit]
Description=Vopak Assistant Daily Harvester

[Service]
Type=oneshot
ExecStart=/usr/bin/node /home/patosoto/geminiclaw/.agent/harvester.ts
StandardOutput=journal
StandardError=journal
